# meta-tizen/meta-tizen-micro/recipes-devtools/nodejs/nodejs/nodejs.service
[Unit]
Description=nodejs service

[Service]

ExecStart=/etc/init.d/nodejs.sh

[Install]
WantedBy=multi-user.target